# systemd/factor-client.service
[Unit]
Description=Factor OctoPrint Client Firmware
After=network.target bluetooth.service
Requires=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
User=factor
Group=factor
WorkingDirectory=/opt/factor-client
ExecStart=/opt/factor-client/venv/bin/python main.py
ExecStartPre=/bin/systemctl start bluetooth
ExecStartPre=/bin/systemctl is-active bluetooth
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target