# systemd/factor-client.service
[Unit]
Description=Factor OctoPrint Client Firmware
After=network.target

[Service]
Type=simple
User=factor
Group=factor
WorkingDirectory=/opt/factor-client
ExecStart=/opt/factor-client/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target